# .ebextensions/supervisor/supervisor_messenger.conf
[program:async-cons]
command=/var/app/current/bin/console messenger:consume async -vvv --time-limit=3600
autostart=true
autorestart=true
user=root
startretries=4
